# etc/skel/.config/autostart/xsnow.desktop
[Desktop Entry]
Exec=xsnow -nomenu
Icon=xsnow
Name=xsnow
Path=
Terminal=False
Type=Application
